# tools/tiktok-monitor/scripts/tiktok-live-probe.service
[Unit]
Description=TikTok live probe (check all accounts)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/eris-TIKTOK
ExecStart=/opt/eris-TIKTOK/.venv/bin/python /opt/eris-TIKTOK/live_worker.py probe
Nice=5

[Install]
WantedBy=multi-user.target
