# tools/tiktok-monitor/scripts/tiktok-account-audit.timer
[Unit]
Description=Run TikTok account ban audit daily at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
